# deploy/server/qiuqiuwriter-backend.service
[Unit]
Description=XingQiu Writer Backend (FastAPI + uvicorn)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/root/xingqiu_writer/backend

EnvironmentFile=/root/xingqiu_writer/backend/.env

# PYTHONPATH 必须指向 src，否则 memos 包找不到
Environment=PYTHONPATH=/root/xingqiu_writer/backend/src

ExecStart=/root/xingqiu_writer/backend/.venv/bin/uvicorn \
    memos.api.ai_api:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=qiuqiuwriter-backend

[Install]
WantedBy=multi-user.target
